# librescoot-modem.service
[Unit]
Description=librescoot modem service
After=network.target

[Service]
Type=simple
WorkingDirectory=/etc/librescoot
ExecStart=/usr/bin/modem-service
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
